# src/python2025_alexigna/hw_05_oop/task2_rpc/worker_supervisor/app/supervisord.conf
[supervisord]
user=root
logfile=/dev/stdout
logfile_maxbytes=0

[program:worker-default]
command=rq worker -P /app/ -c worker default
process_name=%(program_name)s-%(process_num)s
numprocs=4
directory=/app
stopsignal=TERM
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
